# rel/sv/bigcouch/run
#!/bin/bash
exec 2>&1
export HOME=/home/bigcouch
exec chpst -u bigcouch /opt/bigcouch/bin/bigcouch
