# rel/sv/dbcore/run
#!/bin/bash
exec 2>&1
export HOME=/home/dbcore
exec chpst -u dbcore /opt/dbcore/bin/dbcore
